# deploy/scripts/aquillm.service
[Unit]
Description=AquiLLM Service with Docker Compose
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStartPre=/usr/bin/docker compose -f ${PROJECT_ROOT}/deploy/compose/production.yml up -d --remove-orphans
ExecStart=${PROJECT_ROOT}/deploy/scripts/healthcheck.sh
ExecStopPost=/usr/bin/docker compose -f ${PROJECT_ROOT}/deploy/compose/production.yml stop
Restart=always

[Install]
WantedBy=multi-user.target